# etc/init/trypurescript.conf
description "Try PureScript server"
author      "Phil Freeman"

start on started mountall
stop on shutdown

respawn
respawn limit 99 5

script

    exec sudo -u purescript /opt/trypurescript/trypurescript 8050 >> /var/log/trypurescript.log 2>&1
end script

post-start script

end script
